# server/systemd/enshrouded-watchdog.service
[Unit]
Description=Enshrouded idle watchdog
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/enshrouded-watchdog
Restart=always

[Install]
WantedBy=multi-user.target
